# services/tts-core.service
[Unit]
Description=TTSCore - Shared TTS model service
After=graphical-session.target
StartLimitIntervalSec=30

[Service]
Type=simple
ExecStart=%h/projects/tts_core/.venv/bin/python3 -m tts_core
Restart=on-failure
RestartSec=5
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=HF_DATASETS_OFFLINE=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
