# services/tts-core-dashboard.service
[Unit]
Description=TTSCore Dashboard - Web management interface
After=graphical-session.target tts-core.service
Wants=tts-core.service
StartLimitIntervalSec=30

[Service]
Type=simple
ExecStart=%h/projects/tts_core/.venv/bin/python3 -m tts_core.dashboard
Restart=on-failure
RestartSec=5
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
